# scheduler/mirror-digest.service
[Unit]
Description=Mirror digest job

[Service]
Type=oneshot
WorkingDirectory=/path/to/mirror
Environment=PYTHONPATH=/path/to/mirror/src
ExecStart=/usr/bin/python3 /path/to/mirror/src/mirror/cli.py digest
