# scheduler/mirror-coach.service
[Unit]
Description=Mirror scheduled coaching insight export

[Service]
Type=oneshot
WorkingDirectory=/path/to/mirror
Environment=PATH=/home/user/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=uv run --project /path/to/mirror mirror coach --save
