# scheduler/mirror-digest.service
[Unit]
Description=Mirror digest job

[Service]
Type=oneshot
WorkingDirectory=/path/to/mirror
ExecStart=/usr/bin/python3 -m uv run --project /path/to/mirror mirror digest
