# scheduler/mirror-coach.service
[Unit]
Description=Mirror scheduled coaching insight export

[Service]
Type=oneshot
WorkingDirectory=/path/to/mirror
ExecStart=/usr/bin/python3 -m uv run --project /path/to/mirror mirror coach --save
