# scheduler/mirror-digest.timer
[Unit]
Description=Run Mirror digest hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
