# scheduler/mirror-digest.service
[Unit]
Description=Mirror digest job

[Service]
Type=oneshot
WorkingDirectory=/path/to/mirror
Environment=PATH=/home/user/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=uv run --project /path/to/mirror mirror digest
